# examples/mtproto.service.d/resource-limits.conf
# Install: mkdir -p /etc/systemd/system/mtproto.service.d
#         cp this file to /etc/systemd/system/mtproto.service.d/resource-limits.conf
#         systemctl daemon-reload && systemctl restart mtproto.service
#
# Purpose: cap MTProxy cgroup so Xray (3x-ui) keeps CPU/RAM headroom on a small VPS.
# Tune MemoryMax/MemoryHigh and CPUQuota to your RAM and traffic.
#
# Observed 2026-04-12 on sofia (~709 MiB RAM, 1 vCPU): mtproto slave RSS ~100 MiB,
# xray ~62 MiB, x-ui ~54 MiB — leave headroom; raise MemoryMax if OOM in journal.
#
# After applying ExecStart tuning, `systemctl restart mtproto` can fail if global
# PIDs are >65535 (mtproxy bug). Reboot once clears it; see mtproto-proxy-exec-args.example.

[Service]
CPUAccounting=yes
# Lower than default 100: MTProxy yields to x-ui under pressure.
CPUWeight=35
# Hard cap: ~half of one vCPU (100% = 1 core). Increase if MT traffic is dominant.
CPUQuota=55%

# Soft throttle then ceiling for the whole service cgroup (master + workers).
MemoryHigh=130M
MemoryMax=220M

Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=7

# Optional: stop thrashing swap for this unit (kernel >= 5.16)
# MemorySwapMax=64M
