# systemd/cock-mtproxy-monitor.service
[Unit]
Description=cock-monitor MTProxy collect and alerts (oneshot)
Documentation=file:///opt/cock-monitor/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
EnvironmentFile=-/etc/cock-monitor.env
WorkingDirectory=/opt/cock-monitor
ExecStart=/usr/bin/python3 -m cock_monitor mtproxy-collect --env-file /etc/cock-monitor.env

[Install]
WantedBy=multi-user.target

