# systemd/cock-monitor-telegram-bot.timer
[Unit]
Description=Periodic Telegram poll for cock-monitor /status

[Timer]
OnBootSec=3min
# Reply latency is at most this interval; increase on very small VPS if desired.
OnUnitActiveSec=3min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
