# systemd/cock-monitor-leak-investigation.timer
[Unit]
Description=Run cock-monitor leak investigation sampler every 60 seconds
Documentation=file:///opt/cock-monitor/README.md

[Timer]
OnBootSec=1min
OnUnitActiveSec=60s
Unit=cock-monitor-leak-investigation.service

[Install]
WantedBy=timers.target
