# systemd/cock-monitor-daily.timer
[Unit]
Description=Daily cock-monitor chart at 00:05

[Timer]
OnCalendar=*-*-* 00:05:00 Europe/Moscow
Persistent=true

[Install]
WantedBy=timers.target
