# systemd/cock-vless-daily.timer
[Unit]
Description=Daily VLESS traffic report at 00:10 MSK

[Timer]
OnCalendar=*-*-* 00:10:00 Europe/Moscow
Persistent=true

[Install]
WantedBy=timers.target
