# systemd/cock-monitor.service
[Unit]
Description=Light conntrack fill check and Telegram alert (oneshot)
Documentation=file:///opt/cock-monitor/README.md

[Service]
Type=oneshot
# Root is typical: /proc/sys reads are world-readable; conntrack -S often needs privileges.
User=root
# Optional: default EnvironmentFile path on the server (adjust or duplicate in override).
EnvironmentFile=-/etc/cock-monitor.env
WorkingDirectory=/opt/cock-monitor
ExecStart=/usr/bin/python3 -m cock_monitor conntrack-check /etc/cock-monitor.env

[Install]
WantedBy=multi-user.target
