# systemd/cock-monitor-shaper.service
[Unit]
Description=cock-monitor CPU shaper tick
Documentation=file:///opt/cock-monitor/README.md

[Service]
Type=oneshot
User=root
EnvironmentFile=-/etc/cock-monitor.env
WorkingDirectory=/opt/cock-monitor
ExecStart=/opt/cock-monitor/.venv/bin/python -m cock_monitor run shaper /etc/cock-monitor.env

[Install]
WantedBy=multi-user.target
