# systemd/cock-monitor-incident-sampler.service
[Unit]
Description=Incident sampler (network/DNS/conntrack snapshot to JSONL)
Documentation=file:///opt/cock-monitor/README.md

[Service]
Type=oneshot
User=root
EnvironmentFile=-/etc/cock-monitor.env
WorkingDirectory=/opt/cock-monitor
ExecStart=/usr/bin/python3 -m cock_monitor.services.incident_sampler /etc/cock-monitor.env

[Install]
WantedBy=multi-user.target
